# deploy/twidgest-backup.timer
[Unit]
Description=Run twidgest-backup nightly at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
